# crontab/crontab
# Collect data iops for zabbix
* * * * * /etc/zabbix/scripts/iostat-collect.sh /tmp/iostat.out 60
